# debian/usr/lib/systemd/system/amidiauto.service
[Unit]
Description=ALSA MIDI autoconnect daemon
After=sound.target

[Service]
ExecStart=/usr/bin/amidiauto
EnvironmentFile=/etc/environment

[Install]
WantedBy=multi-user.target

